# systemd/plot-butler.service.d/priority.conf
# Dashboard + plot rsync are best-effort relative to recompute farming path.
[Service]
Nice=10
CPUWeight=50
IOWeight=50
OOMScoreAdjust=200
Restart=always
RestartSec=5
