# plot-butler.service
[Unit]
Description=The Plot Butler dashboard and transfer scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=smokey
WorkingDirectory=/home/smokey/plot-butler
EnvironmentFile=-/home/smokey/plot-butler/plot-butler.env
ExecStart=/usr/bin/python3 /home/smokey/plot-butler/plot_butler.py
Restart=always
RestartSec=5
Environment=PLOT_BUTLER_PORT=8088
# Farming-protect knobs (optional overrides):
# Environment=PLOT_BUTLER_BWLIMIT_KBPS=12000
# Environment=PLOT_BUTLER_HARVESTER_PAUSE_S=15
# Environment=PLOT_BUTLER_API_TOKEN=
# Environment=PLOT_BUTLER_BIND=0.0.0.0

[Install]
WantedBy=multi-user.target

