# systemd/chia-recompute.service.d/limits.conf
[Service]
# Many concurrent farmer connections under load
LimitNOFILE=65536
TasksMax=infinity
