# systemd/chia-recompute.service.d/priority.conf
# Farming recompute is latency-critical. Prefer it over plot shipping / dashboard.
[Service]
Nice=-10
CPUWeight=500
IOWeight=500
OOMScoreAdjust=-500
# Keep restart aggressive after transient GPU/driver blips
Restart=always
RestartSec=5
TimeoutStopSec=30
